# supervisor/conf.d/supervisor.conf
[supervisord]
nodaemon=true

[program:apache2]
user=ids
command=/bin/bash -c "source /etc/apache2/envvars && export APACHE_RUN_USER=ids && exec /usr/sbin/apache2 -DFOREGROUND"
